# bootc/rootfs/etc/systemd/user/signal-cli.service
[Unit]
Description=signal-cli JSON-RPC daemon
Documentation=https://github.com/AsamK/signal-cli
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/signal-cli --output=json daemon --http=0.0.0.0:8181 --receive-mode=on-start
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
